# rootfs/etc/services.d/samba/run
#!/usr/bin/with-contenv sh
exec /usr/sbin/samba --model=$PROCESS_MODEL -i </dev/null